# roles/elasticsearch/templates/cassandra.ini
[program:cassandra]
command=env PATH=$PATH:/bin:/usr/bin/:/usr/java/jdk1.8.0_45/bin/ JAVA_HOME=/usr/java/jdk1.8.0_45 /data/bigdata/app/cassandra/bin/cassandra -f
user=fangyu.fu
autostart=true
autorestart=true
startsecs=10
startretries=999
log_stdout=true
log_stderr=true
logfile=/data/var/log/supervisor/cassandra.out
logfile_maxbytes=20MB
logfile_backups=10
